# deploy/stock-chatbot.service
[Unit]
Description=Sanjuk Stock Chatbot (Telegram)
After=network.target

[Service]
Type=simple
User=kanzaka110
WorkingDirectory=/home/kanzaka110/Sanjuk-Stock-Simulator
ExecStart=/bin/bash -c 'export $(grep -v "^#" /home/kanzaka110/Sanjuk-Stock-Simulator/.env | xargs) && /home/kanzaka110/Sanjuk-Stock-Simulator/venv/bin/python -u main.py chatbot'
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
