# deploy/stock-bot.service
[Unit]
Description=Sanjuk Stock Bot + Market Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kanzaka110
WorkingDirectory=/home/kanzaka110/Sanjuk-Stock-Simulator
EnvironmentFile=/home/kanzaka110/Sanjuk-Stock-Simulator/.env
ExecStart=/bin/bash -c 'cd /home/kanzaka110/Sanjuk-Stock-Simulator && /home/kanzaka110/Sanjuk-Stock-Simulator/venv/bin/python -u main.py bot'
Restart=always
RestartSec=15

# 로그 — journalctl -u stock-bot -f 로 확인
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
